# services/lifeos-server.service
[Unit]
Description=LifeOS Backend Server (LifeOnline)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/xionglei/LifeOnline/LifeOS/packages/server
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/home/linuxbrew/.linuxbrew/sbin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/home/xionglei/LifeOnline/LifeOS/packages/server/.env.production
ExecStart=/home/linuxbrew/.linuxbrew/bin/pnpm start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lifeos-server

[Install]
WantedBy=default.target
